# launchd/com.shoirhi.dotfiles.autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shoirhi.dotfiles.autoupdate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/dotfiles/scripts/auto-update.sh</string>
  </array>

  <!-- 毎日 12:00 に実行。実行時刻にスリープ/電源オフでも、次回起動・復帰時に一度だけ実行される -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/.local/state/dotfiles/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.local/state/dotfiles/launchd.err.log</string>
</dict>
</plist>
